# src/etc/systemd/system/om-telegram-alert.service
[Unit]
Description=OpenWeatherMap Telegram Alert

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/om-telegram-alert.py
